# .config/systemd/user/file_sync.timer
[Unit]
Description=File sync service to save certain files

[Service]
Type=simple
ExecStart=

[Install]
WantedBy=multi-user.target
