# .config/systemd/user/task-warrior-remind.service
[Unit]
Description=Task Warrior Reminder Service

[Service]
ExecStart=%h/.scripts/task-warrior-remind
