# .config/systemd/user/xscreensaver.service
[Unit]
Description=XScreenSaver service

[Service]
Type=simple
ExecStart=xscreensaver

[Install]
WantedBy=default.target
